# files/ansible-pull.timer
[Unit]
Description=Run ansible-pull

[Timer]
OnCalendar=daily
Persistent=true 

[Install]
WantedBy=timers.target
